# files/usr/lib/systemd/user/dotfile-manager.service
[Unit]
Description=Manage default dotfiles

[Service]
Type=oneshot
ExecStart=/usr/libexec/dotfile-manager

[Install]
WantedBy=default.target
