# files/usr/lib/systemd/user/mango-session.target
[Unit]
Description=MangoWC Session Target
Requires=graphical-session.target
After=graphical-session.target
