# files/usr/lib/systemd/user/user-up.timer
[Unit]
Description=Upgrade system 1h after boot

[Timer]
OnBootSec=1h

[Install]
WantedBy=timers.target
