# files/usr/lib/systemd/system/prefix-dedupe.service
[Unit]
Description=Deduplicate common wine prefixes weekly
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/libexec/prefix-dedupe
