# files/usr/etc/profile.d/disable-nix.sh
#!/bin/sh
alias nix='echo "Did you mean nx?"'
